# packaging/claude-discord-bridge.service
[Unit]
Description=Claude Code <-> Discord bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/claude-discord-bridge serve
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
